# rootfs/etc/s6-overlay/s6-rc.d/thelounge/run
#!/usr/bin/with-contenv sh

# run (from /usr/local/bin)
exec thelounge --config "host=0.0.0.0" start
